# pkg/ds/tpl/50-gizmo.rules.tpl
KERNEL=="sd*", ATTRS{vendor}=="RPI", ATTRS{model}=="RP2", RUN+="/usr/local/bin/gss-loader"
